# neo-tracker.cron
# do update/current at 10am, 12pm, 2pm, 4pm, sleeping a random amount before doing so.
# at 6pm, do the daily update.
0 10,12,14,16 * * 1-5 nobody    perl -le 'sleep rand 3600' && curl --silent --insecure --request POST https://localhost/app/update/current >/dev/null
0 18 * * 1-5          nobody    curl --silent --insecure --request POST https://localhost/app/update/daily >/dev/null
